# root/etc/s6-overlay/s6-rc.d/init-mod-universal-tor/run
#!/usr/bin/with-contenv bash

echo "tor" >> /mod-repo-packages-to-install.list