# root/etc/s6-overlay/s6-rc.d/svc-mod-universal-tor/run
#!/usr/bin/with-contenv bash
if [[ ${TOR} -eq 1 ]]; then
    exec s6-setuidgid tor tor -f ${TOR_CONFIG:-/etc/tor/torrc}
else
    exec bash -c -- 'sleep 15; while true; do sleep infinity; done;'
fi
